# system/etc/profile.d/rbenv.sh
export RBENV_ROOT=/opt/github/rbenv
export PATH="${RBENV_ROOT}/bin:${PATH}"
eval "$(rbenv init --no-rehash -)"
